# run-minutely.cron
#!/bin/bash

cd  /users/a/r/areagan/fun/twitter/scraping

python scrape-oauth2.py

